# service/log/run
#!/bin/sh
     exec setuidgid root multilog t ./main
